# backend/deploy/payday-worker.service
# Systemd Service - Celery Worker
# 用于生产环境管理异步任务处理
#
# 使用方法:
# 1. 修改 ExecStart 中的路径为实际部署路径
# 2. 复制到 /etc/systemd/system/payday-worker.service
# 3. 重载 systemd: systemctl daemon-reload
# 4. 启用服务: systemctl enable payday-worker
# 5. 启动服务: systemctl start payday-worker
# 6. 查看状态: systemctl status payday-worker
# 7. 查看日志: journalctl -u payday-worker -f

[Unit]
Description=Celery Worker for PayDay Backend
After=network.target redis.service mysql.service
Requires=redis.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/path/to/payday/backend
Environment="PATH=/path/to/venv/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/path/to/venv/bin/celery -A app.celery_app worker \
    --pidfile=/var/run/celery/payday-worker.pid \
    --logfile=/var/log/celery/payday-worker.log \
    --loglevel=info \
    --concurrency=4
ExecStop=/path/to/venv/bin/celery -A app.celery_app multi stopwait \
    --pidfile=/var/run/celery/payday-worker.pid \
    --logfile=/var/log/celery/payday-worker.log
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=payday-worker

[Install]
WantedBy=multi-user.target
